# playground-api/playground-api.service
[Unit]
Description=Zen-C Playground API
After=network.target

[Service]
User=root
WorkingDirectory=/opt/zenc-playground-api
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
